# crontabs/daily-cron.sh
#!/bin/bash

cd /var/SaintObjetBot

php daily-cron.php bluesky
php daily-cron.php mastodon
