# crontabs/daily-cron.sh
#!/bin/bash

cd /home/tobozo/src/SaintObjetBot

/usr/bin/php daily-cron.php bluesky
/usr/bin/php daily-cron.php mastodon
