# crontabs/daily-crawl.sh
#!/bin/bash

cd /var/SaintObjetBot

php keyword-liker.php bluesky
php keyword-liker.php mastodon
