# crontabs/weekly-stats.sh
#!/bin/bash

cd /var/SaintObjetBot

php bsky-stats.php
php mstdn-stats.php update

