# crontabs/yearly-calendar.sh
#!/bin/bash

cd /var/SaintObjetBot

php calendar.php $(date +%Y)
